# Code Bundle-20250227T232753Z-001/Code Bundle/Chapter 08/my_service.service
[Unit]
Description=My Custom Service
After=network.target
[Service]
ExecStart=/usr/bin/my_custom_service
Restart=always
[Install]
WantedBy=multi-user.target
